# example/templates/supervisord.conf
[supervisord]
logfile=/dev/null
pidfile=/var/run/supervisord.pid
nodaemon=true

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0777
chown=nobody:nogroup

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
unix:///var/run/supervisor.sock

{{ if getenv "ETCD" }}
[program:confd]
user=root
command=confd -node {{ getenv "ETCD" }} -config-file /app/confd.toml
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
{{ end }}

[program:healthcheck]
command=/app/bin/healthcheck
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
user=app

[program:example]
command=gunicorn -b 0.0.0.0:8080 app:app
directory=/app/example
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
user=app
